# root/etc/s6-overlay/s6-rc.d/svc-netbox/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

exec \
    s6-notifyoncheck -d -n 30 -w 5000 \
        cd /app/netbox/netbox s6-setuidgid abc /usr/sbin/uwsgi uwsgi.ini
